# vnc-app/supervisor/app.conf
[program:xeyes]
environment=DISPLAY=:1
command=xeyes