# app/supervisor/app.conf
[program:gnome-calculator]
environment=DISPLAY=:1
command=gnome-calculator